# systemd/speechd-ng.service
[Unit]
Description=Speech Server Daemon Next-Gen
Documentation=https://github.com/StarTuz/speechd-ng
After=graphical-session.target

[Service]
Type=simple
ExecStart=%h/.local/bin/speechd-ng
Restart=on-failure
BusName=org.speech.Service
Environment=RUST_LOG=info

# ========================
# SECURITY HARDENING
# ========================

# Filesystem Protection
ProtectSystem=full
ProtectHome=false
PrivateTmp=true

# Allow access to models (some are in cache)
# BindReadOnlyPaths can be used for extra safety if ProtectHome=true, 
# but for user service we'll stay simple for now.

# Prevent privilege escalation
NoNewPrivileges=true

# Protect kernel and system configuration
ProtectKernelTunables=true
ProtectKernelModules=true
ProtectKernelLogs=true
ProtectControlGroups=true
ProtectClock=true
ProtectHostname=true

# Device access (needed for audio and potentially GPU/Screenshot)
PrivateDevices=false

# Network: Allow localhost (Ollama) and Hubs
RestrictAddressFamilies=AF_UNIX AF_INET AF_INET6

# Misc
UMask=0077

# Misc
UMask=0077

[Install]
WantedBy=default.target
